# build/etc/supervisor/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisor.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/lock/supervisord.pid
nodaemon=true
minfds=1024 
minprocs=200
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf

[program:sshd]
command=/opt/base/bin/supervisor-sshd.sh
process_name=%(program_name)s
numprocs=1
directory=/root
priority=100
autostart=true
startsecs=0
startretries=3
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/sshd.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=1
stderr_logfile=/var/log/supervisor/sshd.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=1
environment=PROC_NAME="%(program_name)s",PROC_NUM=%(process_num)d

[program:rclone_mount]
command=/opt/base/bin/supervisor-rclone.sh
process_name=="%(program_name)s_%(process_num)s"
numprocs=%(ENV_RCLONE_MOUNT_COUNT)s
directory=/mnt
priority=200
autostart=true
startsecs=10
startretries=3
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/rclone-mount-%(process_num)s.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=1
stderr_logfile=/var/log/supervisor/rclone-mount-%(process_num)s.error.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=1
environment=PROC_NAME="%(program_name)s",PROC_NUM=%(process_num)d


; Tail all the logs to stdout/err for 'docker logs' commmand
; Use this process for capture & events
[program:logtail]
command=/opt/base/bin/logtail.sh
process_name=%(program_name)s
numprocs=1
directory=/var/log/supervisor
; All other programs should be up first
priority=10000
autostart=true
startsecs=0
startretries=3
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
; Main process stdout
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_capture_maxbytes=1MB
; Main process stderr
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stderr_capture_maxbytes=1MB
environment=PROC_NAME="%(program_name)s"
