# build/COPY_ROOT/etc/supervisor/supervisord/conf.d/caddy.conf
[program:caddy]
command=supervisor-caddy.sh
process_name=%(program_name)s
numprocs=1
directory=%(ENV_WORKSPACE)s
priority=50000
autostart=true
startsecs=3
startretries=3
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/caddy.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=1
stderr_logfile=/var/log/supervisor/caddy.error.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=1
environment=PROC_NAME="%(program_name)s"
