# build/COPY_ROOT/etc/supervisor/supervisord/conf.d/logtail.conf
; Tail all the logs to stdout/err for 'docker logs' commmand
; Use this process for capture & events
[program:logtail]
command=supervisor-logtail.sh
process_name=%(program_name)s
numprocs=1
directory=/var/log/supervisor
; All other programs should be up first
priority=100000
autostart=true
startsecs=0
startretries=3
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
; Main process stdout
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_capture_maxbytes=1MB
; Main process stderr
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stderr_capture_maxbytes=1MB
environment=PROC_NAME="%(program_name)s"
