# build/COPY_ROOT/etc/supervisor/supervisord/conf.d/redirector.conf
[program:redirector]
command=supervisor-redirector.sh
process_name=%(program_name)s
numprocs=1
directory=/opt/ai-dock/fastapi/redirector
priority=100
autostart=true
startsecs=0
startretries=3
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/redirector.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=1
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
stderr_logfile_backups=0
